# Dockerfile/zabbix-xxl/container-files-zabbix/etc/supervisor.d/zabbix-proxy.conf
[program:zabbix-proxy]
command = sudo -u %(ENV_ZP_User)s /usr/local/sbin/zabbix_proxy -f -c /usr/local/etc/zabbix_proxy.conf
autorestart = true
