# Dockerfile/zabbix-xxl/container-files-zabbix/etc/supervisor.d/zabbix-agent.conf
[program:zabbix-agent]
command = /config/ds.sh %(ENV_ZA_PidFile)s sudo -u %(ENV_ZA_User)s /usr/local/sbin/zabbix_agentd -f -c /usr/local/etc/zabbix_agentd.conf
autorestart = true
