# Dockerfile/dockbix-xxl-3.4/container-files-zabbix/etc/supervisor.d/snmptrapd.conf
[program:snmptrapd]
command = /usr/sbin/snmptrapd -f -Lo -n -M+/etc/zabbix/snmp/mibs -m+ALL
autorestart = true
auto_start = true

stderr_logfile = NONE
stdout_logfile = NONE
