# Dockerfile/zabbix-3.0/container-files-zabbix/etc/supervisor.d/zabbix-agent.conf
[program:zabbix-agent]
command = /config/ds.sh /tmp/zabbix_agentd.pid sudo -u zabbix /usr/local/sbin/zabbix_agentd -f -c /usr/local/etc/zabbix_agentd.conf
autorestart = true
