# Dockerfile/zabbix-3.0/container-files-zabbix/etc/supervisor.d/zabbix-server.conf
[program:zabbix-server]
command = /config/ds.sh %(ENV_ZS_PidFile)s sudo -u %(ENV_ZS_User)s /usr/local/sbin/zabbix_server -f -c /usr/local/etc/zabbix_server.conf
autorestart = true
