# Dockerfile/zabbix-2.4/container-files-base/etc/supervisor.d/php-fpm.conf
[program:php-fpm]
command = php-fpm
autorestart = true
stdout_logfile = /data/logs/php-fpm.log
stderr_logfile = /data/logs/php-fpm.log
