# Dockerfile/zabbix-server-2.4/container-files-zabbix/etc/supervisor.d/snmpd.conf
[program:snmptrapd]
command = /usr/sbin/snmptrapd -f -Lo -On -M+/etc/zabbix/snmp/mibs -m+ALL
autorestart = true

#[program:snmptt]
#command = /usr/sbin/snmptt --daemon --debug=2
#autorestart = true
