# Dockerfile/zabbix-server-2.4/container-files-zabbix/etc/supervisor.d/snmptrapd.conf
[program:snmptrapd]
command = /usr/sbin/snmptrapd -f -Lo -n -M+/etc/zabbix/snmp/mibs -m+ALL
#command = /usr/sbin/snmptrapd -Lsd -f
autorestart = true
